# conf/mongodb.conf
[program:mongodb]
command=/usr/local/bin/gunicorn wsgi -c /root/forum/conf/gunicorn.config.py
directory=/root/forum
autostart=true
autorestart=true
